# crontab
30 3 * * * cd /CROUStillantBackup && uv run /CROUStillantBackup/__main__.py >> /CROUStillantBackup/crontab.log
